# webapps/files/default/package.compare.cron
*/5 * * * * root /usr/bin/package.compare.sh >> /var/log/www-pull.log 2>&1 /var/log/www-pull.log 
