# rsyslog/files/default/syslog-emails.cron
50 06 * * * root /usr/local/bin/syslog-email.sh 2>&1 | logger -t syslog-emails
